# debian/sogo.cron.daily
#!/bin/sh

SOGOSPOOL=/var/spool/sogo

/usr/sbin/tmpreaper 24 "$SOGOSPOOL"
find "$SOGOSPOOL" -type d -empty -exec /bin/rmdir -p {} \; 2> /dev/null
