# debian/sogo.postrm
#!/bin/sh
set -e
userdel sogo||true
update-rc.d SOGo remove ||true
